# app-emulation/lxd/files/lxd.initd.1
#!/sbin/openrc-run
# Copyright 2004-2018 Gentoo Foundation, Copryight 2018 Funtoo Solutions
# Distributed under the terms of the GNU General Public License v2

DAEMON=/usr/libexec/lxd
PIDFILE=/run/lxd.pid

extra_commands="stopall"

depend() {
	need net
	use lxcfs
}

start() {
	# required for lxd in lxd:
	modprobe -q ip6_tables
	ebegin "Starting lxd service"
	LD_LIBRARY_PATH=/usr/__LIBDIR__/lxd start-stop-daemon --start \
		--pidfile ${PIDFILE} \
		--exec ${DAEMON} \
		--background \
		--make-pidfile \
		-- \
		${LXD_OPTIONS}
	eend $?
}

stop() {
	if [ "$RC_GOINGDOWN" = "YES" ] || [ "$RC_REBOOT" = "YES" ]; then
		stopall
	else
		ebegin "Stopping lxd service (but not containers)"
		start-stop-daemon --stop --quiet -R TERM/45 -p ${PIDFILE}
		eend $?
	fi
}

stopall() {
	ebegin "Stopping lxd service and containers"
	if "${DAEMON}" shutdown; then
		/etc/init.d/lxd zap
		rm -f ${PIDFILE}
		fi
	eend $?
}
