# sys-power/powerbus/files/powerbus.initd
#!/sbin/openrc-run

name="funtoo-powerbus"
description="Funtoo Power Management Bus"
command=/usr/bin/funtoo-powerbus
command_user="root:root"
command_background=false
pidfile=/run/funtoo/powerbus.pid
supervisor=supervise-daemon

start_pre() {
	install -d /run/funtoo -m 755 -o root -g root
}
